# ajax-gitmobile/ajax-relay/deploy/ajax-relay.service
[Unit]
Description=ajax-relay (EC2 wake/sleep + SSH-over-WS bridge)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/ajax-relay
EnvironmentFile=/etc/ajax-relay/relay.env
ExecStart=/usr/bin/java -jar /opt/ajax-relay/relay-0.1.0.jar
Restart=on-failure
RestartSec=5
SuccessExitStatus=143
StandardOutput=append:/var/log/ajax-relay/relay.log
StandardError=append:/var/log/ajax-relay/relay.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/ajax-relay /var/log/ajax-relay

[Install]
WantedBy=multi-user.target
